# tools/repo_mgmt/sync.cron
### Sync from the OBS repository per 4 hours.
01 */4 * * * /var/www/sync/script/hw-mainline-std.sh
15 */4 * * * /var/www/sync/script/hw-pandora-xfce4.sh
30 */4 * * * /var/www/sync/script/hw-mainline-adv.sh
45 */4 * * * /var/www/sync/script/oerv-factory.sh
59 */4 * * * /var/www/sync/script/oerv-2203.sh
