# tools/repo_mgmt/sync.cron
### Sync from the OBS repository per 4 hours.
01 */4 * * * /var/www/sync/script/hw-mainline-std.sh
30 */4 * * * /var/www/sync/script/hw-mainline-adv.sh
